# config/supervisor/react.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:adb]
directory=/code
command=adb start-server
autostart=true
autorestart=false
startsecs=0
stdout_logfile=/var/log/react.out
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/react.out
stderr_logfile_maxbytes=0

[program:react]
directory=/code
command=react-native run-android
autostart=true
autorestart=true
stdout_logfile=/var/log/react.out
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/react.out
stderr_logfile_maxbytes=0
